# stargate.service
# /etc/systemd/system/stargate.service

[Unit]
Description=Stargate Composition
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=no
RestartSec=5
WorkingDirectory=/etc/stargate
EnvironmentFile=/etc/stargate/.env
ExecStartPre=/usr/bin/docker-compose pull
ExecStart=/usr/bin/docker-compose up --remove-orphans
ExecReload=/usr/bin/docker-compose up -d --remove-orphans
ExecStop=/usr/bin/docker-compose down -v

[Install]
WantedBy=multi-user.target
